# bin/crontab/collect
#!/bin/sh

directory="$( cd `dirname $0` && pwd )/../.."

(
    flock -n 1 || exit 1

    php ${directory}/bin/console app:collect --env=prod

) 1>~/var/tmp/`basename $0`