# bin/crontab/filter
#!/bin/sh

directory="$( cd `dirname $0` && pwd )/../.."

(
    flock -n 2 || exit 1

    php ${directory}/bin/console app:filter --env=prod

) 2>~/var/tmp/`basename $0`